# templates/etc/profile.d/env.sh
GITHUB_PAT=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
